# crontab
@reboot /home/hp/.config/misc/randommac.sh
@reboot /usr/bin/intel-undervolt daemon
